# packaging/deb/postinst
#!/bin/bash
set -e

if ! id -u corrijameptbr >/dev/null 2>&1; then
  adduser --system --group --home /var/lib/corrija_me_pt_br corrijameptbr
fi

mkdir -p /var/lib/corrija_me_pt_br
chown -R corrijameptbr:corrijameptbr /var/lib/corrija_me_pt_br

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload || true
  systemctl enable corrija-me-pt-br.service || true
  systemctl restart corrija-me-pt-br.service || true
fi

exit 0
